# luci/applications/luci-app-tfortis-teleport/root/etc/uci-defaults/80_tfortis-teleport
#!/bin/sh

touch /etc/config/tfortis-teleport
uci set tfortis-teleport.@tfortis-teleport[0]=tfortis-teleport
uci set tfortis-teleport.@tfortis-teleport[0].state='1'
uci commit

return 0
