# luci/applications/luci-app-tfortis-io/root/etc/uci-defaults/tfortis-io
config tfortis-io
	option	tamper_state 1
	option out1_state 'enable'
