# luci/applications/luci-app-tfortis-ups/root/etc/uci-defaults/80_tfortis-ups
#!/bin/sh

touch /etc/config/tfortis-ups
uci set tfortis-ups.first=first
uci set tfortis-ups.second=second
uci commit

return 0
